# spec/comodule/deployment/test/platform/experiment/config/crontab/make_cache.txt.erb
0 20 * * * /bin/bash -l -c 'cd /home/ec2-user/projects/<%= config.app_name %> && /root/.rbenv/shims/ruby bin/rails runner MakeCache.update -e production'
